# root/etc/s6-overlay/s6-rc.d/adsbhub/finish
#!/command/with-contenv bash

# exit container stop s6
kill 1
