# root/etc/services.d/fr24feed/run
#!/bin/execlineb -P
/fr24feed/fr24feed_amd64/fr24feed