# root/etc/services.d/fr24feed/run
#!/usr/bin/with-contenv bash
if [ "$SERVICE_ENABLE_FR24FEED" != "false" ]
then
	/bin/decoration -c blue -p "[fr24feed] " -e "/fr24feed/fr24feed_amd64/fr24feed"
else
	tail -f /dev/null
fi