# root/etc/services.d/piaware/run
#!/usr/bin/with-contenv bash
if [ "$SERVICE_ENABLE_PIAWARE" != "false" ]
then
	/bin/decoration -c green -p "[piaware] " -e "/usr/bin/piaware" -a "-plainlog"
else
	tail -f /dev/null
fi
