# root/etc/services.d/piaware/run
#!/usr/bin/with-contenv bash
if [ "$SERVICE_ENABLE_PIAWARE" != "false" ]
then
	/usr/bin/piaware -plainlog
else
	tail -f /dev/null
fi
