# root/etc/s6-overlay/s6-rc.d/confd-fr24feed/up
#!/command/with-contenv bash

/etc/s6-overlay/s6-rc.d/confd-fr24feed/script
